# clashd.service
[Unit]
Description=clash-verge-service as sudo

[Service]
ExecStart=sudo clash-verge-service

[Install]
WantedBy=multi-user.target

